# systemd/dvr.service
[Unit]
Description=DVR HDMI Capture Application (Weston kiosk)
After=tc358743.service network.target
Wants=tc358743.service

[Service]
User=pi
Group=pi
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/usr/bin/weston --shell=kiosk-shell.so
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
