# systemd/tc358743.service
[Unit]
Description=TC358743 HDMI-CSI2 capture device init (EDID + format)
After=dev-video0.device
Wants=dev-video0.device
Before=dvr.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/tc358743-init.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
