# pi-gen/stage-dvr/03-provision/files/dvr-wifi-restore.service
[Unit]
Description=Restore saved Wi-Fi connections onto read-only root
DefaultDependencies=no
After=local-fs.target
Before=NetworkManager.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/dvr-wifi-restore.sh

[Install]
WantedBy=NetworkManager.service
