# pi-gen/stage-dvr/03-provision/files/dvr-provision.service
[Unit]
Description=DVR first-boot provisioning (enable read-only root)
After=local-fs.target
Before=getty@tty1.service
ConditionPathExists=!/boot/firmware/.dvr-provisioned

[Service]
Type=oneshot
ExecStart=/usr/local/bin/dvr-provision.sh
StandardOutput=tty
TTYPath=/dev/tty1

[Install]
WantedBy=multi-user.target
