# pi-gen/stage-dvr/03-provision/files/dvr-provision.service
[Unit]
Description=DVR first-boot provisioning (enable read-only root)
After=multi-user.target
ConditionPathExists=!/boot/firmware/.dvr-provisioned

[Service]
Type=oneshot
ExecStart=/usr/local/bin/dvr-provision.sh

[Install]
WantedBy=multi-user.target
